# crontab
# OPIM - Open Platform Infrastructure Monitoring running each 3 minutes
*/3 * * * *	/opt/telefonica/opim/opim upload >> /var/log/opim.log 2>&1
